# build/yunsh-headtracking.service
[Unit]
Description=YUNSH OS Head Tracking Daemon (3DoF IMU)
Documentation=https://github.com/ljcccc999/yunsh-os
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/yunsh-headtracking
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
