# build/yunsh-headtracking.service
[Unit]
Description=YUNSH OS Head Tracking
After=multi-user.target
[Service]
Type=simple
ExecStart=/usr/bin/yunsh-headtracking
Restart=always
[Install]
WantedBy=multi-user.target
